# data/systemd/system/evdbindkeyd.service.in
[Unit]
Description=evdbindkeyd service

[Service]
Type=simple
ExecStart=@bindir@/evdbindkeyd

[Install]
WantedBy=multi-user.target
